# rfremix-config-18/yum-makecache.cron
#! /bin/bash

LOGFILE=/var/log/yum-makecache.cron
STAMP=`date +%c`:
echo $STAMP >>$LOGFILE
/usr/bin/yum clean all >>$LOGFILE
/usr/bin/yum makecache >>$LOGFILE
